# deploy/axiom-daily-review.service
[Unit]
Description=Axiom daily review snapshot

[Service]
Type=oneshot
WorkingDirectory=/opt/axiom
EnvironmentFile=/opt/axiom/.env
ExecStart=/opt/axiom/.venv/bin/python /opt/axiom/scripts/run_logged_automation.py --job-id review_day --root /opt/axiom --days-offset -1 --utc-offset +08:00
