# deploy/axiom-backup.timer
[Unit]
Description=Run Axiom daily backup

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=10m
Unit=axiom-backup.service

[Install]
WantedBy=timers.target
