# deploy/axiom-daily-inbox-action.service
[Unit]
Description=Axiom daily inbox action dry-run snapshot

[Service]
Type=oneshot
WorkingDirectory=/opt/axiom
EnvironmentFile=/opt/axiom/.env
ExecStart=/opt/axiom/.venv/bin/python /opt/axiom/scripts/save_inbox_action_snapshot.py --root /opt/axiom --days-offset -1 --utc-offset +08:00 --force
