# deploy/axiom-inbox-processing.service
[Unit]
Description=Axiom inbox processing snapshot

[Service]
Type=oneshot
WorkingDirectory=/opt/axiom
EnvironmentFile=/opt/axiom/.env
ExecStart=/opt/axiom/.venv/bin/python /opt/axiom/scripts/run_logged_automation.py --job-id inbox_report --root /opt/axiom --utc-offset +08:00
