# deploy/axiom-daily-inbox-action-history.service
[Unit]
Description=Axiom daily inbox action history snapshot

[Service]
Type=oneshot
WorkingDirectory=/opt/axiom
EnvironmentFile=/opt/axiom/.env
ExecStart=/opt/axiom/.venv/bin/python /opt/axiom/scripts/save_inbox_action_history_snapshot.py --root /opt/axiom --window day --days-offset -1 --utc-offset +08:00 --details --force
