# deploy/axiom-weekly-inbox-action-history.service
[Unit]
Description=Axiom weekly inbox action history snapshot

[Service]
Type=oneshot
WorkingDirectory=/opt/axiom
EnvironmentFile=/opt/axiom/.env
ExecStart=/opt/axiom/.venv/bin/python /opt/axiom/scripts/run_logged_automation.py --job-id inbox_action_history_week --root /opt/axiom --days-offset -1 --utc-offset +08:00
