# deploy/axiom-daily-audio-transcribe.service
[Unit]
Description=Axiom daily audio transcribe snapshot

[Service]
Type=oneshot
WorkingDirectory=/opt/axiom
EnvironmentFile=/opt/axiom/.env
ExecStart=/opt/axiom/.venv/bin/python /opt/axiom/scripts/run_logged_automation.py --job-id audio_transcribe_day --root /opt/axiom --days-offset -1 --utc-offset +08:00 --skip-when-unavailable
